# mta/etc/cron.d/postmap
* * * * * root postmap /opt/holonet/postfix_domains
* * * * * root postmap /opt/holonet/postfix_lmtp
